# deploy/rocket.socket
[Unit]
Description=rocket socket

[Socket]
ListenStream=/run/rocket.sock

[Install]
WantedBy=multi-user.target